# res/package/scripts/postinstall.sh
#!/bin/sh

systemctl daemon-reload
systemctl enable prometheus-nginxlog-exporter